# test/assets/test_pkg/test_pkg.service
[Unit]
Description=test pkg

[Service]
ExecStart=/usr/bin/test_pkg

[Install]
WantedBy=multi-user.target